# systemd/proxmox_agent.timer
[Unit]
Description=Run Sensor Reader Every 5 Minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
Unit=proxmox_agent.service

[Install]
WantedBy=timers.target